# etc/systemd/system/flask.service
[Unit]

Description=Flask_Book

[Service]

User=jenkins

ExecStart=/var/lib/jenkins/workspace/example_freestyle/script/installation.sh

[Install]

WantedBy=multi-user.target